# university/cfg/systemd/gpio.timer
[Unit]
Description=Try to get GPIO state every 3 seconds

[Timer]
OnBootSec=3sec
OnUnitActiveSec=3sec

[Install]
WantedBy=timers.target
